# com.christian.ai-agent-team.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.christian.ai-agent-team</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/christian/Repos/My-AI-Agent-Team/start-permanent-dev.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/christian/Repos/My-AI-Agent-Team</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/christian/Repos/My-AI-Agent-Team/logs/startup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/christian/Repos/My-AI-Agent-Team/logs/startup.error.log</string>
</dict>
</plist>
